# templates/suricata.service.j2
[Unit]
Description=Suricata IDS/IPS
After=network.target

[Service]
ExecStart=/usr/bin/suricata -c /etc/suricata/suricata.yaml -i {{ interface }}

[Install]
WantedBy=default.target